# ops/systemd/hud-mcp.service
# hud-mcp.service — HUD MCP daemon (HTTP/SSE, bearer auth, ACL)
#
# IMPORTANT: This unit runs the long-lived mcp-hud HTTP/SSE daemon.
# It is NOT a replacement for Emily's per-session stdio MCP — that path
# continues to work unaffected. This unit serves external callers
# (Hermes gateway, Hermes MacBook) over Tailscale on port 7610.
#
# Deploy steps (run as root on the Hetzner server):
#   cp ops/systemd/hud-mcp.service /etc/systemd/system/hud-mcp.service
#   systemctl daemon-reload
#   systemctl enable --now hud-mcp.service
#   journalctl -u hud-mcp -f   # verify clean start
#
# The daemon reads secrets from /srv/hud/secrets/mcp.env (mode 600,
# owner hud). That file is NEVER committed to git.
# See ops/secrets/mcp-tokens.example.yaml and ops/secrets/mcp-acl.example.yaml
# for the required structure.
#
# Token rotation: edit /srv/hud/secrets/mcp-tokens.yaml, then
#   systemctl restart hud-mcp.service
# (The daemon loads tokens at startup; a full restart is required for rotation.
#  Since the StreamableHTTPServerTransport is stateless, restarts are safe:
#  no in-flight SSE connections are dropped — each request is self-contained.)

[Unit]
Description=HUD MCP Daemon (HTTP/SSE transport, bearer auth, ACL)
Documentation=https://github.com/kevinaton/hud
After=network.target
Wants=network.target

[Service]
Type=simple
User=hud
Group=hud
Slice=hud.slice
WorkingDirectory=/srv/hud/app

# Environment — real secrets live here (mode 600, owner hud)
EnvironmentFile=/srv/hud/secrets/mcp.env

# Main process: Node 22 running the compiled mcp-hud daemon entry point.
# MCP_TRANSPORT and PORT are set in mcp.env or overridden here for clarity.
# NOTE: dist-new/ is the hud-owned build output (the original dist/ was built by
# root during Ticket 35/36 and is read-only). dist-new/ is rebuilt by `pnpm --filter mcp-hud build`
# after root removes the old dist/ (see Ticket 37 deploy notes).
ExecStart=/usr/bin/node packages/mcp-hud/dist-new/index.js

# Restart policy — survive transient errors and Hermes restarts
Restart=on-failure
RestartSec=5s
StartLimitIntervalSec=120s
StartLimitBurst=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hud-mcp

# Hardening — matches hud-web.service pattern; ReadWritePaths limited to data
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/srv/hud/data

# Additional hardening
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
LockPersonality=true
# MemoryDenyWriteExecute intentionally omitted: Node.js JIT requires
# writable+executable memory pages. Enabling this would prevent startup.
RestrictRealtime=true
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
